# scripts/dev/com.rbtechbot.bde-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- DEV ONLY — NOT bundled in the BDE app DMG.
     Before use, replace all occurrences of /Users/RBTECHBOT/Documents/Repositories/BDE
     with your actual repo path, and /Users/RBTECHBOT with your home directory.
     See scripts/dev/README.md for full setup instructions. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.rbtechbot.bde-watcher</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/RBTECHBOT/Documents/Repositories/BDE/scripts/dev/bde-watcher.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/RBTECHBOT/Documents/Repositories/BDE</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/RBTECHBOT</string>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/bde-watcher.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/bde-watcher.log</string>
</dict>
</plist>
